# linux/deb/postinst
#!/bin/bash
set -e

APP_DIR=/opt/seestar_alp
SEESTAR_USER=seestar
CONFIG_DIR=/etc/seestar

# Create dedicated system user if it doesn't exist
if ! id -u "$SEESTAR_USER" >/dev/null 2>&1; then
    useradd --system --no-create-home --shell /usr/sbin/nologin \
        --home-dir "$APP_DIR" "$SEESTAR_USER"
fi

# Create log directory writable by the service user
mkdir -p "$APP_DIR/logs"
chown "$SEESTAR_USER:$SEESTAR_USER" "$APP_DIR/logs"

# Create config directory
mkdir -p "$CONFIG_DIR"

# On first install, generate config.toml from the bundled example.
# On upgrade, the existing config is left untouched.
if [ ! -f "$CONFIG_DIR/config.toml" ]; then
    cp "$APP_DIR/device/config.toml.example" "$CONFIG_DIR/config.toml"
    # Default to listening on all interfaces (not just loopback)
    sed -i 's/127\.0\.0\.1/0.0.0.0/g' "$CONFIG_DIR/config.toml"
    sed -i 's|^log_prefix =.*|log_prefix = "logs/"|g' "$CONFIG_DIR/config.toml"
fi

# Symlink the config into the location the app expects.
# Force-create so upgrades replace any stale regular file the app may have
# created (config.py copies the example if the target is missing and writable).
ln -sf "$CONFIG_DIR/config.toml" "$APP_DIR/device/config.toml"

# The venv and Python runtime are pre-built into the package.
# Just ensure the seestar user owns them.
chown -R "$SEESTAR_USER:$SEESTAR_USER" "$APP_DIR/.venv" "$APP_DIR/.python"

# Apply sysctl settings (IPv6 disable)
if command -v sysctl >/dev/null 2>&1; then
    sysctl --system >/dev/null 2>&1 || true
fi

# Enable and start the main seestar service
systemctl daemon-reload
systemctl enable seestar.service
systemctl restart seestar.service

HOSTNAME=$(hostname -I 2>/dev/null | awk '{print $1}' || hostname)
echo ""
echo "|---------------------------------------------|"
echo "| Seestar ALP installation complete            |"
echo "|                                             |"
printf "| Access at: http://%-26s|\n" "${HOSTNAME}:5432"
echo "|                                             |"
echo "| Logs:   journalctl -u seestar               |"
echo "| Status: systemctl status seestar            |"
echo "| Config: /etc/seestar/config.toml            |"
echo "|                                             |"
echo "| Optional INDI support:                      |"
echo "|  Install the seestar-alp-indi package       |"
echo "|---------------------------------------------|"
