# linux/deb/postinst
#!/bin/bash
set -e

APP_DIR=/opt/seestar_alp
SEESTAR_USER=seestar
CONFIG_DIR=/etc/seestar

# Create dedicated system user if it doesn't exist
if ! id -u "$SEESTAR_USER" >/dev/null 2>&1; then
    useradd --system --no-create-home --shell /usr/sbin/nologin \
        --home-dir "$APP_DIR" "$SEESTAR_USER"
fi

# Create log directory writable by the service user
mkdir -p "$APP_DIR/logs"
chown "$SEESTAR_USER:$SEESTAR_USER" "$APP_DIR/logs"

# Create config directory
mkdir -p "$CONFIG_DIR"

# On first install, generate config.toml from the bundled example.
# On upgrade, the existing config is left untouched.
if [ ! -f "$CONFIG_DIR/config.toml" ]; then
    cp "$APP_DIR/device/config.toml.example" "$CONFIG_DIR/config.toml"
    # Default to listening on all interfaces (not just loopback)
    sed -i 's/127\.0\.0\.1/0.0.0.0/g' "$CONFIG_DIR/config.toml"
    sed -i 's|^log_prefix =.*|log_prefix = "logs/"|g' "$CONFIG_DIR/config.toml"
fi

# Symlink the config into the location the app expects.
# Force-create so upgrades replace any stale regular file the app may have
# created (config.py copies the example if the target is missing and writable).
ln -sf "$CONFIG_DIR/config.toml" "$APP_DIR/device/config.toml"

# ---------------------------------------------------------------------------
# Python virtual environment
#
# Packages built with Docker buildx ship a pre-built .venv (and the
# uv-managed Python in .python).  In that case just fix ownership and skip
# the slow runtime install.
#
# Packages built without Docker ship a bundled uv binary and build the venv
# at install time.  This path also handles manual dpkg installs on the target.
# ---------------------------------------------------------------------------
mkdir -p "$APP_DIR/.local" "$APP_DIR/.cache"

if [ -d "$APP_DIR/.venv" ]; then
    echo "Using pre-built Python virtual environment..."
    chown -R "$SEESTAR_USER:$SEESTAR_USER" "$APP_DIR/.venv"
else
    # Resolve uv — prefer bundled binary, fall back to system uv
    if [ -x "$APP_DIR/.local/bin/uv" ]; then
        UV="$APP_DIR/.local/bin/uv"
    elif command -v uv >/dev/null 2>&1; then
        UV=uv
    else
        echo "Error: uv not found and no pre-built venv is present." >&2
        echo "Re-install the package built with Docker buildx, or install uv first." >&2
        exit 1
    fi

    echo "Installing Python dependencies (this may take a few minutes on first install)..."
    mkdir -p "$APP_DIR/.venv"
    chown -R "$SEESTAR_USER:$SEESTAR_USER" "$APP_DIR/.venv" "$APP_DIR/.local"
    HOME="$APP_DIR" su -s /bin/sh "$SEESTAR_USER" -c "
        '$UV' venv --python 3.13 '$APP_DIR/.venv' &&
        '$UV' pip install --python '$APP_DIR/.venv/bin/python' \
            -r '$APP_DIR/requirements.txt'
    " 2>&1 | tail -10
fi

chown -R "$SEESTAR_USER:$SEESTAR_USER" "$APP_DIR/.local" "$APP_DIR/.cache"

# Apply sysctl settings (IPv6 disable)
if command -v sysctl >/dev/null 2>&1; then
    sysctl --system >/dev/null 2>&1 || true
fi

# Enable and start the main seestar service
systemctl daemon-reload
systemctl enable seestar.service
systemctl restart seestar.service

HOSTNAME=$(hostname -I 2>/dev/null | awk '{print $1}' || hostname)
echo ""
echo "|---------------------------------------------|"
echo "| Seestar ALP installation complete            |"
echo "|                                             |"
printf "| Access at: http://%-26s|\n" "${HOSTNAME}:5432"
echo "|                                             |"
echo "| Logs:   journalctl -u seestar               |"
echo "| Status: systemctl status seestar            |"
echo "| Config: /etc/seestar/config.toml            |"
echo "|                                             |"
echo "| Optional INDI support:                      |"
echo "|  Install the seestar-alp-indi package       |"
echo "|---------------------------------------------|"
